# etc/cron.weekly/certbot
#!/bin/bash

/snap/bin/certbot renew -q
