# etc/cron.weekly/certbot
#!/bin/bash

# Renewing certificates

/snap/bin/certbot renew -q
